# rootfs/etc/init.d/kill-frecon
#!/sbin/openrc-run
# OpenRC init script for kill-frecon
# Kills ChromeOS frecon processes that interfere with shimboot

description="Kill ChromeOS frecon processes"

depend() {
    need localmount
    after bootmisc
}

start() {
    ebegin "Killing ChromeOS frecon processes"
    
    # Kill frecon-lite if running
    if pgrep -x frecon-lite > /dev/null 2>&1; then
        pkill -9 frecon-lite 2>/dev/null || true
    fi
    
    # Kill any remaining frecon processes
    pkill -9 frecon 2>/dev/null || true
    
    eend 0
}