# rootfs/etc/init.d/kill-frecon
#!/sbin/openrc-run
# OpenRC init script for kill-frecon
# Kills ChromeOS frecon processes that interfere with shimboot/Xorg

description="Release ChromeOS frecon console for Xorg"
command="/usr/local/bin/kill_frecon"
command_args=""

depend() {
    need localmount
    after bootmisc modules
    before xdm display-manager
}
