# rootfs/etc/init.d/kill-frecon
#!/sbin/openrc-run
# OpenRC helper for manually releasing ChromeOS frecon.
# Do not add this service to boot; xdm runs /usr/local/bin/kill_frecon right
# before LightDM starts so the screen is blank for the shortest possible time.

description="Release ChromeOS frecon console for Xorg"
command="/usr/local/bin/kill_frecon"
command_args=""

depend() {
    need localmount
    after bootmisc modules
    before xdm
}
