# rootfs/etc/init.d/mdev
#!/sbin/openrc-run

# Alpine-compatible mdev service for shimboot Gentoo.

description="Populate /dev with busybox mdev"
command="/usr/local/bin/shimboot-mdev-scan"

depend() {
    need devfs sysfs
    before udev udev-trigger hwdrivers modules shimboot-hwmods shimboot-xfce
}
