# rootfs/etc/init.d/shimboot-hwmods
#!/sbin/openrc-run

description="Load Chromebook hardware/input modules for shimboot"
command="/usr/local/bin/shimboot-load-hwmods"

depend() {
    need sysfs
    after udev-trigger modules
    before display-manager xdm shimboot-xfce
}
